# systemd/abemadl.timer
[Unit]
Description=Daily Timer for AbemaTV Video Donwloader

[Timer]
OnCalendar=*-*-* 4:20:00
Unit=abemadl.service

[Install]
WantedBy=timers.target
